# supervisor.cfg
[program:DTKLP-numbers-recongnition]
process_name=DTKLP-numbers-recongnition
command=sh scrips/runinenv.sh dependencies/venv python3 main.py
directory=~/api-rest-car-number/
user=user
autostart=true
autorestart=true
redirect_stderr = True
stdout_logfile= /var/log/%(program_name)s/%(process_num)s-out.log
stderr_logfile= /var/log/%(program_name)s/%(process_num)s-err.log